# launchd/com.sanosano.kimochimax.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sanosano.kimochimax.health</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/sanosano/projects/kimochi-max/health_monitor.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/sanosano/projects/kimochi-max</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/kimochi_health.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/kimochi_health.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
